# SSE_Python_supervisord.conf
; Ideally, if you run default-install.sh, you shouldn't have to modify anything here.

[program:SSE_Python]
command=python2.7 /opt/sse/main.py
process_name=SSE_Python
autostart=true
autorestart=true
startsecs=60
startretries=3
exitcodes=0,1
stopsignal=KILL

[supervisord]
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
user=root
